# target/linux/mediatek/filogic/base-files/etc/uci-defaults/04_led_migration
. /lib/functions/migrations.sh

board=$(board_name)

case "$board" in
tplink,archer-ax80-v1)
	migrate_leds ':chan=-'
	;;
esac

remove_devicename_leds

migrations_apply system

exit 0
